# arch/debian/postrm.in
#!/bin/sh
#*=====================================================================*/
#*    serrano/prgm/project/hop/2.2.x/arch/debian/debian/postrm.in      */
#*    -------------------------------------------------------------    */
#*    Author      :  Manuel Serrano                                    */
#*    Creation    :  Tue Dec 18 10:55:54 2007                          */
#*    Last change :  Tue Oct 25 08:05:32 2011 (serrano)                */
#*    Copyright   :  2007-11 Manuel Serrano                            */
#*    -------------------------------------------------------------    */
#*    Remove the HOP desktop entries from the configuration files      */
#*=====================================================================*/

#*---------------------------------------------------------------------*/
#*    maemo                                                            */
#*---------------------------------------------------------------------*/
if [ "@DEBIAN@ " = "maemo " ]; then
  path="/etc/hildon-desktop"
  conf=home.conf
  home_data="hop.desktop"

  sed "/$home_data/d" $path/$conf > home_conf.tmp
  mv home_conf.tmp $path/$conf
fi

#*---------------------------------------------------------------------*/
#*    init script                                                      */
#*---------------------------------------------------------------------*/
if [ "$1" = purge ]; then
  update-rc.d hop remove
fi
