# arch/archlinux/systemd/hop.socket.in
[Unit]
Description=Hop Web Broker Sockets

[Socket]
ListenStream=/var/run/hop/socket
ListenStream=8080

[Install]
WantedBy=sockets.target
