# arch/archlinux/systemd/hop.service.in
[Unit]
Description=Hop Web Broker

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/conf.d/hop
ExecStart=${DAEMON} -p ${PORT} --eval "(hop-user-set! \"${DAEMONUSER}\")" --acept-kill ${DAEMONARGS}
Execstop=${DAEMON} -p ${PORT} --eval "(hop-user-set! \"${DAEMONUSER}\")" --kill ${DAEMONARGS}

[Install]
WantedBy=multi-user.target
